# deployment/systemd/beauty-landing.service
[Unit]
Description=Beauty Platform - Landing frontend (prod)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/beauty-dev/apps/landing
Environment=NODE_ENV=production
Environment=PORT=7004
EnvironmentFile=/opt/beauty-dev/.env.production
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beauty-landing

[Install]
WantedBy=multi-user.target
