# deployment/systemd/beauty-payment-service.service
[Unit]
Description=Beauty Platform - Payment Service (prod)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/beauty-dev/services/payment-service
Environment=NODE_ENV=production
Environment=PORT=7029
EnvironmentFile=/opt/beauty-dev/.env.production
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=10s
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beauty-payment-service

[Install]
WantedBy=multi-user.target
